# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=/app/backend/server.js

[program:nginx]
command=nginx -c /app/docker/nginx.conf
