# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:node]
command=/app/backend/server.js

[program:nginx]
command=nginx -g "daemon off;"
