# docker/supervisord.conf
[supervisord]
nodaemon=true
user=myuser

[program:node]
command=/app/backend/server.js

[program:nginx]
command=nginx -g "daemon off;http.server.listen ${PORT:-80};"
